# projects/nig/backend/cron/analysis_management.cron
SHELL=/bin/bash

0 */8 * * * set -a && source /etc/rapydo-environment && /usr/bin/python3 /code/nig/scripts/init_pipeline.py >> /logs/init_pipeline.log
